# supervisord.conf
[program:mongod]
command=mongod --port 27018
autostart=true
autorestart=true
